# php-fpm/crontab
* * * * * sh -c 'php /var/www/apisiaa/artisan schedule:run >> /dev/null  2>&1'
